# backend/railway/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgresql]
command=/usr/bin/postgres -D /var/lib/postgresql/data
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql.err

[program:redis]
command=/usr/bin/redis-server --daemonize no
autostart=true
autorestart=true
priority=2
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.err

[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
autostart=true
autorestart=true
priority=3
stdout_logfile=/var/log/supervisor/rabbitmq.log
stderr_logfile=/var/log/supervisor/rabbitmq.err

[program:auth-service]
command=node /app/services/auth-service/dist/main.js
autostart=true
autorestart=true
priority=10
environment=PORT=3001,NODE_ENV=production
stdout_logfile=/var/log/supervisor/auth.log
stderr_logfile=/var/log/supervisor/auth.err

[program:wallet-service]
command=node /app/services/wallet-service/dist/main.js
autostart=true
autorestart=true
priority=10
environment=PORT=3002,NODE_ENV=production
stdout_logfile=/var/log/supervisor/wallet.log
stderr_logfile=/var/log/supervisor/wallet.err

[program:escrow-service]
command=node /app/services/escrow-service/dist/main.js
autostart=true
autorestart=true
priority=10
environment=PORT=3003,NODE_ENV=production
stdout_logfile=/var/log/supervisor/escrow.log
stderr_logfile=/var/log/supervisor/escrow.err

[program:delivery-service]
command=node /app/services/delivery-service/dist/main.js
autostart=true
autorestart=true
priority=10
environment=PORT=3004,NODE_ENV=production
stdout_logfile=/var/log/supervisor/delivery.log
stderr_logfile=/var/log/supervisor/delivery.err

[program:dispute-service]
command=node /app/services/dispute-service/dist/main.js
autostart=true
autorestart=true
priority=10
environment=PORT=3005,NODE_ENV=production
stdout_logfile=/var/log/supervisor/dispute.log
stderr_logfile=/var/log/supervisor/dispute.err

[program:notification-service]
command=node /app/services/notification-service/dist/main.js
autostart=true
autorestart=true
priority=10
environment=PORT=3006,NODE_ENV=production
stdout_logfile=/var/log/supervisor/notification.log
stderr_logfile=/var/log/supervisor/notification.err

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err
