# backend/railway/supervisord-render.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=1
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis]
command=/usr/bin/redis-server --daemonize no --loglevel warning
autostart=true
autorestart=true
priority=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:auth-service]
command=sh -c 'sleep 2 && node /app/services/auth-service/dist/main.js'
autostart=true
autorestart=true
priority=10
startretries=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:wallet-service]
command=sh -c 'sleep 4 && node /app/services/wallet-service/dist/main.js'
autostart=true
autorestart=true
priority=10
startretries=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:escrow-service]
command=sh -c 'sleep 4 && node /app/services/escrow-service/dist/main.js'
autostart=true
autorestart=true
priority=10
startretries=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:delivery-service]
command=sh -c 'sleep 5 && node /app/services/delivery-service/dist/main.js'
autostart=true
autorestart=true
priority=10
startretries=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dispute-service]
command=sh -c 'sleep 5 && node /app/services/dispute-service/dist/main.js'
autostart=true
autorestart=true
priority=10
startretries=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:notification-service]
command=sh -c 'sleep 5 && node /app/services/notification-service/dist/main.js'
autostart=true
autorestart=true
priority=10
startretries=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
