# backend/railway/supervisord-render.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
pidfile=/var/run/supervisord.pid

[program:redis]
command=/usr/bin/redis-server --daemonize no --loglevel warning
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:auth-service]
command=node /app/services/auth-service/dist/main.js
autostart=true
autorestart=true
priority=10
environment=PORT=3001,NODE_ENV=production,DATABASE_URL="%(ENV_DATABASE_URL)s"
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:wallet-service]
command=node /app/services/wallet-service/dist/main.js
autostart=true
autorestart=true
priority=10
environment=PORT=3002,NODE_ENV=production,DATABASE_URL="%(ENV_DATABASE_URL)s"
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:escrow-service]
command=node /app/services/escrow-service/dist/main.js
autostart=true
autorestart=true
priority=10
environment=PORT=3003,NODE_ENV=production,DATABASE_URL="%(ENV_DATABASE_URL)s"
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:delivery-service]
command=node /app/services/delivery-service/dist/main.js
autostart=true
autorestart=true
priority=10
environment=PORT=3004,NODE_ENV=production,DATABASE_URL="%(ENV_DATABASE_URL)s"
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:dispute-service]
command=node /app/services/dispute-service/dist/main.js
autostart=true
autorestart=true
priority=10
environment=PORT=3005,NODE_ENV=production,DATABASE_URL="%(ENV_DATABASE_URL)s"
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:notification-service]
command=node /app/services/notification-service/dist/main.js
autostart=true
autorestart=true
priority=10
environment=PORT=3006,NODE_ENV=production,DATABASE_URL="%(ENV_DATABASE_URL)s"
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/null
stderr_logfile=/dev/null
